# packaging/ubuntu/postinst
echo "Starting agent"
if [ $(cat /proc/1/comm) = init ]
then
    start amazon-ssm-agent || true
elif [ $(cat /proc/1/comm) = systemd ]
then
    systemctl enable amazon-ssm-agent
    systemctl start amazon-ssm-agent
    systemctl daemon-reload
fi
